# server/benchmark-server.service
[Unit]
Description=benchmark server service

[Service]
ExecStart=/root/benchmark-japan-vps
Restart=always
RestartSec=60
StandardOutput=append:/root/benchmark-server.log
StandardError=append:/root/benchmark-server-err.log

[Install]
WantedBy=default.target
