# server/benchmark-server.service
[Unit]
Description=benchmark server service

[Service]
WorkingDirectory=/root
ExecStart=/root/benchmark-japan-vps
Restart=always
RestartSec=60
StandardOutput=append:/root/logs/server.log
StandardError=append:/root/logs/err.log

[Install]
WantedBy=default.target
